# supervisord.conf
[supervisord]

nodaemon=true

[program:rabbitmq]

command=/start.sh
